# Apps/GanchoMenuBarHelper/com.johnny4young.gancho.menubar-helper.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- Active LaunchAgent for SMAppService.agent. BundleProgram is the helper
	     embedded in the app; its exact-owner watchdog self-exits when Gancho is
	     gone, so RunAtLoad cannot leave an orphaned menu-bar process. -->
	<key>Label</key>
	<string>com.johnny4young.gancho.menubar-helper.agent</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/GanchoMenuBarHelper</string>
	<key>RunAtLoad</key>
	<true/>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>com.johnny4young.gancho</string>
	</array>
</dict>
</plist>
